# deploy/systemd/autosales-api.service
[Unit]
Description=[YourBrand] Auto Sales API
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/opt/AutoAI_Sales
Environment="PATH=/opt/AutoAI_Sales/.venv/bin"
ExecStart=/opt/AutoAI_Sales/scripts/start_api.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
